# motion.conf
[program:motion]
command=/usr/bin/motion
autostart=false
autorestart=false
user=motion

