# deployment/mosquitto.conf
[program:mosquitto]
command=/usr/sbin/mosquitto
autostart=true
autorestart=true
startretries=1000
user=mosquitto
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
