# deployment/mosquitto.conf
[program:mosquitto]
command=/usr/sbin/mosquitto
autostart=true
autorestart=true
startretries=1000
user=mosquitto
redirect_stderr=true
stdout_logfile=/var/log/mosquitto.log
