# deployment/motion.conf
[program:motion]
command=/usr/bin/motion
autostart=false
autorestart=false
user=motion
stdout_logfile=/var/log/motion.log
redirect_stderr=true
