# deployment/motion.conf
[program:motion]
command=/usr/bin/motion
autostart=false
autorestart=false
user=motion
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
