# deployment/nfc_scan.conf
[program:nfc_scan]
command=/usr/local/bin/things2c nfc_scan
autostart=true
autorestart=true
startretries=1000
stdout_logfile=/var/log/nfc_scan.log
redirect_stderr=true