# deployment/nfc_scan.conf
[program:nfc_scan]
command=/usr/local/bin/things2c nfc_scan
autostart=true
autorestart=true
stderr_logfile=/var/log/nfc_scan.err.log
stdout_logfile=/var/log/nfc_scan.out.log
