# deployment/watchdog.conf
[program:watchdog]
command=/usr/local/bin/things2c watchdog
autostart=true
autorestart=true
startretries=1000
stdout_logfile=/var/log/watchdog.log
redirect_stderr=true
